# RL/data/clawgym_train/task_1795/input_files/input/systemd/press-release-fetch.service
[Unit]
Description=Fetch official press releases mentioning Alessandro Amitrano
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/utente/newsroom
ExecStart=/usr/bin/python3 /home/utente/newsroom/scripts/fetch_amitranonews.py --config /home/utente/newsroom/config/fetcher.yaml
User=utente

[Install]
WantedBy=multi-user.target
