# RL/data/clawgym_train/task_1795/input_files/input/systemd/press-release-fetch.timer
[Unit]
Description=Run press-release-fetch daily at 09:00

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
Unit=press-release-fetch.service

[Install]
WantedBy=timers.target
